# systemd/db_backup@.service
[Unit]
Description=DB backup: %I

[Service]
User=postgres
ExecStart=/bin/bash /usr/local/bin/db_backup.sh %i
