# systemd/db_backup@.timer
[Unit]
Description=Daily backup of %I

[Timer]
OnCalendar=daily
Unit=db_backup@%i.service

[Install]
WantedBy=multi-user.target
